# config/photurisiii-ocr.service.d/environment.conf
[Service]
Environment="AZUREKEY=xxx"
Environment="AZUREURL=xxx"
Environment="STRIPE_SECRET_KEY=xxx"
Environment="STORAGE_BUCKET=xxx"
Environment="STRIPE_WEB_HOOK_KEY=xxx"